# deploy/vroumi.service
# Service systemd pour faire tourner Vroumi en permanence sur un Raspberry Pi.
#
# Installation :
#   1. Adaptez WorkingDirectory, User et le chemin de npm si besoin (voir
#      docs/RASPBERRY-PI.md).
#   2. sudo cp deploy/vroumi.service /etc/systemd/system/vroumi.service
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable --now vroumi
#
# Suivi :     sudo systemctl status vroumi
# Journaux :  journalctl -u vroumi -f

[Unit]
Description=Vroumi - suivi d'entretien automobile
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Utilisateur qui possède le dossier du projet (souvent "pi").
User=pi
Group=pi

# Chemin ABSOLU du projet sur le Raspberry Pi.
WorkingDirectory=/home/pi/vroumi

# Variables d'environnement (chargées depuis le fichier .env du projet).
EnvironmentFile=/home/pi/vroumi/.env
# Écoute sur toutes les interfaces réseau (accessible depuis le LAN).
# PORT 3001 pour coexister avec Oudiral (qui occupe déjà le 3000) sur le même Pi.
# Si Vroumi est seul sur la machine, vous pouvez repasser à 3000.
Environment=NODE_ENV=production
Environment=HOSTNAME=0.0.0.0
Environment=PORT=3001

# Démarrage du serveur de production Next.js.
# Vérifiez le chemin de npm avec `which npm` (souvent /usr/bin/npm).
ExecStart=/usr/bin/npm run start

# Redémarrage automatique en cas de plantage.
Restart=on-failure
RestartSec=5

# Robustesse / sécurité de base.
NoNewPrivileges=true
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
